# conf_files/bitmonero.conf
start on runlevel [2345]
stop on shutdown

expect daemon

#setuid bob
exec /home/bob/bitmonero/build/release/bin/bitmonerod --detach --rpc-bind-ip 192.168.1.8 --show-time-stats 0

