# conf_files/cpu_miner.conf
start on started poolnode
stop on stopping poolnode

console log
chdir /monerodo/yam/
exec ./yamgeneric -c yam.cfg
