# conf_files/mos_monerodoip.conf
start on started network-services
stop on shutdown

exec /monerodo/change_ip.sh
