# conf_files/mos_monerodoip.conf
start on runlevel [2345]
stop on shutdown

exec /monerodo/change_ip.sh
