# conf_files/mossys_monerodoip.conf
start on (net-device-up IFACE!=lo) and runlevel [2345]
stop on shutdown

exec /monerodo/change_ip.sh
