# conf_files/upstart/mos_mininodo.conf
start on nodowallet_sync
stop on stopping mos_nodowallet 
console log 
respawn 
respawn limit 10 10 
post-start exec sleep 5 
chdir /monerodo/MiniNodo/ 
exec nodejs MiniNodo.js -p poop 

