# conf_files/upstart/mos_monerowallet.conf
start on bitmonero_sync
stop on stopping mos_bitmonero
console log
respawn
respawn limit 10 10
pre-start exec logrotate -f /etc/logrotate.d/upstart
post-start exec sh -c 'tail -n +0 --pid=$$ -f /var/log/upstart/mos_monerowallet.log | { sed "/Run net_service/ q" && kill $$ ;}'

